# .config/autostart/startup-commands.desktop
[Desktop Entry]
Name=Startup commands
Type=Application
NoDisplay=true
Exec=/usr/lib/lightdm/guest-session-auto.sh
